# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground /opt/zookeeper/conf/zoo_sample.cfg
autostart=true
autorestart=true
